# deploy/systemd/cloudflared-wsiannotate.service
[Unit]
Description=Cloudflare Tunnel for WSI Annotate
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared --config /etc/cloudflared/config.yml tunnel run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
